# studio-ops/systemd/achilles-queue-runner.service
# achilles-queue-runner.service — one poll cycle of the ops-task-queue local runner.
#
# The runner is the PULL side of the cloud-clock/local-runner design (validated
# 2026-07-02): claim due/overdue rows from the service-role-only ops_task_queue,
# run the named ops-agents in DRY-RUN, write bookkeeping + a device_heartbeats
# row. Offline catch-up needs no logic here — overdue rows wait in the queue.
# Guarded writes stay with the main session: the runner never passes --apply.
#
# Install: bin/studio-runner-install (copies units, creates user+env, enables timer).

[Unit]
Description=ACHILLES ops-task-queue runner (one poll cycle, dry-run agents only)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=achilles-runner
Group=achilles-runner
WorkingDirectory=/studio/clients/_ops-agents
# 0600 root-owned env file; holds SUPABASE_URL / SUPABASE_SERVICE_ROLE_KEY /
# OPS_DEVICE_ID by value — the ONLY place they live for the runner.
EnvironmentFile=/etc/achilles/queue-runner.env
ExecStart=/usr/bin/env npx tsx agents/queue-runner.ts
TimeoutStartSec=15min

# Hardening (the runner device is "keys to the kingdom" — validation action item 3)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/studio/clients/_ops-agents
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
