# studio-ops/systemd/achilles-queue-runner-docker@.timer
# achilles-queue-runner-docker@.timer — boot +2min, every 15min, Persistent.
# %i = container name; pairs with achilles-queue-runner-docker@.service.

[Unit]
Description=Timer for ACHILLES queue runner (docker exec into %i)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
RandomizedDelaySec=30
Unit=achilles-queue-runner-docker@%i.service

[Install]
WantedBy=timers.target
